# context-manager/systemd/dcm-dashboard.service
[Unit]
Description=DCM Dashboard (Next.js)
PartOf=dcm.target
After=dcm-api.service
Wants=dcm-api.service
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
Type=simple
WorkingDirectory=__DASHBOARD_ROOT__
Environment=NODE_ENV=production
Environment=PORT=__DASHBOARD_PORT__
Environment=PATH=__SERVICE_PATH__
Environment=NEXT_PUBLIC_API_URL=http://127.0.0.1:__API_PORT__
Environment=NEXT_PUBLIC_WS_URL=ws://127.0.0.1:__WS_PORT__
ExecStartPre=/usr/bin/npx next build
ExecStart=/usr/bin/npx next start -p __DASHBOARD_PORT__
KillMode=mixed
KillSignal=SIGTERM
TimeoutStartSec=120
TimeoutStopSec=15
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dcm-dashboard
LimitNOFILE=65536
MemoryMax=4G

[Install]
WantedBy=dcm.target
